# rootfs/etc/s6-overlay/s6-rc.d/cups-server/finish
#!/usr/bin/with-contenv bashio
set -euo pipefail

/usr/bin/with-contenv /bin/bash /usr/share/cups-state-backup.sh
